# deploy/blocweather-frontend.service
[Unit]
Description=BlocWeather Frontend
After=network.target blocweather-backend.service

[Service]
Type=simple
User=blocweather
WorkingDirectory=/opt/blocweather/frontend
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node /opt/blocweather/frontend/build/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
